# com.ytmp3.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ytmp3.api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/caris/.gemini/antigravity/scratch/yt-mp3-api/run.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/caris/.gemini/antigravity/scratch/yt-mp3-api/logs/stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/caris/.gemini/antigravity/scratch/yt-mp3-api/logs/stdout.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/caris/.gemini/antigravity/scratch/yt-mp3-api</string>
</dict>
</plist>
